# define/supervisor/conf.d/php-fpm.conf
[supervisord]
nodaemon=true

[program:nginx]
nodaemon=true
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:php-fpm]
nodaemon=true
command = /usr/local/sbin/php-fpm
autostart=true
autorestart=true